# src/JT808.DotNetty.Tests/JT808.DotNetty.Hosting/JT808.DotNetty.Hosting.ini
[program:JT808.DotNetty.Hosting]
command=dotnet JT808.DotNetty.Hosting.dll
directory=/data/dotnetty_test
user=root
stdout_logfile=/data/supervisordlog/JT808.DotNetty.Hosting.log
redirect_stderr=true